# packaging/ubuntu/postinst
echo "Starting agent"
if [ $(cat /proc/1/comm) = init ]
then
    start qvagent || true
elif [ $(cat /proc/1/comm) = systemd ]
then
    systemctl enable qvagent
    systemctl start qvagent
    systemctl daemon-reload
fi
