# packaging/ubuntu/prerm
echo "Stopping agent"
if [ $(cat /proc/1/comm) = init ]
then
    stop qvagentt || true
elif [ $(cat /proc/1/comm) = systemd ]
then
    systemctl stop qvagentt
    systemctl disable qvagentt
    systemctl daemon-reload
fi
